# testdata/golden/minimal/runtime/ownbase-auth.container
# Generated by OwnBase — do not hand-edit.
# Changes here will be overwritten and flagged as a tamper signal.

[Unit]
Description=OwnBase container: ownbase-auth
After=network-online.target

[Container]
Image=localhost/ownbase-auth:local
ContainerName=ownbase-auth
Network=ownbase-auth-net.network
Network=ownbase-internal.network
Volume=ownbase-auth-data:/data
PublishPort=127.0.0.1:41000:8080
# BuildSource=services/auth
# BuildRef=v2.1.0
NoNewPrivileges=true
DropCapability=ALL

[Service]
Restart=always
TimeoutStartSec=30

[Install]
WantedBy=default.target
