# testdata/golden/minimal/runtime/ownbase-hello.container
# Generated by OwnBase — do not hand-edit.
# Changes here will be overwritten and flagged as a tamper signal.

[Unit]
Description=OwnBase container: ownbase-hello
After=network-online.target

[Container]
Image=localhost/ownbase-hello:local
ContainerName=ownbase-hello
Network=ownbase-auth-net.network
Network=ownbase-hello-net.network
Network=ownbase-internal.network
Volume=ownbase-hello-data:/data
PublishPort=127.0.0.1:3000:3000
# Requires=auth
# BuildSource=apps/hello
NoNewPrivileges=true
DropCapability=ALL

[Service]
Restart=always
TimeoutStartSec=30

[Install]
WantedBy=default.target
