# testdata/golden/dev-tls/runtime/ownbase-app.container
# Generated by OwnBase — do not hand-edit.
# Changes here will be overwritten and flagged as a tamper signal.

[Unit]
Description=OwnBase container: ownbase-app
After=network-online.target

[Container]
Image=localhost/ownbase-app:local
ContainerName=ownbase-app
Network=ownbase-app-net.network
Network=ownbase-internal.network
Volume=ownbase-app-data:/data
PublishPort=127.0.0.1:3000:3000
# BuildSource=apps/app
NoNewPrivileges=true
DropCapability=ALL

[Service]
Restart=always
TimeoutStartSec=30

[Install]
WantedBy=default.target
